# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/call-chain.conf
[Service]
# Front-door: rules first, model only if ambiguous (1s wall)
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER_TIMEOUT=2.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
# Clear single tools skip ornith cold start
Environment=AIPC_DAILY_DIRECT_TOOLS=1
# Voice chat: skip mem0 recall for first token; remember still async
Environment=AIPC_VOICE_MEM0_RECALL=0
Environment=AIPC_MEM0_ASYNC=1
# Long jobs: heartbeat + progress trail
Environment=AIPC_TASK_ASYNC=1
Environment=AIPC_TASK_HEARTBEAT_S=4
Environment=AIPC_TASK_PROGRESS_MAX=40
# Fail-fast short LLM
Environment=AIPC_LLM_VOICE_TIMEOUT=6
Environment=AIPC_LLM_VOICE_MAX_RETRIES=0
