# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/classifier.conf
[Service]
# Model-first daily routing (resident-small). Rules only for greet/farewell.
# AIPC_CLASSIFIER=rules was a speed hack that forced keyword-only daily —
# product default is always (model classify under timeout wall).
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER_TIMEOUT=2.0
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
