# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/speed.conf
[Service]
# Model-first classify (always). Rules still short-circuit greet/farewell.
# Timeout wall keeps free-chat from hanging if classifier model is cold.
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER_TIMEOUT=2.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
Environment=AIPC_LLM_VOICE_TIMEOUT=30
Environment=AIPC_LLM_VOICE_MAX_RETRIES=0
Environment=AIPC_LLM_REQUEST_TIMEOUT=45
Environment=AIPC_LLM_MAX_RETRIES=0
Environment=AIPC_MEM0_VOICE_TIMEOUT=0.8
Environment=AIPC_MEM0_ASYNC=1
Environment=AIPC_VOICE_MEM0=1
Environment=AIPC_VOICE_MEM0_RECALL=0
Environment=AIPC_DAILY_DIRECT_TOOLS=1
Environment=AIPC_DAILY_LLM_TIMEOUT=45
Environment=AIPC_CLARIFY=auto
Environment=AIPC_TASK_ASYNC=1
