# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zz-call-chain.conf
[Service]
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_TIMEOUT=2.5
Environment=AIPC_LLM_VOICE_TIMEOUT=30
Environment=AIPC_LLM_REQUEST_TIMEOUT=45
Environment=AIPC_DAILY_DIRECT_TOOLS=1
Environment=AIPC_VOICE_MEM0_RECALL=0
