# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-uncensored-agent.conf
[Service]
# Model-judged policy (no per-topic keyword gates in code):
#   chat  → assistant-gemma (uncensored balanced; model decides how to answer)
#   classify → qwythos-9b (fast; routes by task shape only, never refuses topics)
#   daily/Hermes tools → coder-agentic (uncensored aggressive)
#   screen → vlm-screen (speed)
# NPU resident-small optional max-speed chat (more refusals):
#   Environment=AIPC_SUPERVISOR_MODEL=resident-small
Environment=AIPC_SUPERVISOR_MODEL=assistant-gemma
Environment=AIPC_CLASSIFIER_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_TIMEOUT=3.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
Environment=AIPC_DAILY_MODEL=coder-agentic
Environment=AIPC_DAILY_DIRECT_TOOLS=1
Environment=AIPC_SCREEN_VLM=vlm-screen
Environment=AIPC_LLM_VOICE_TIMEOUT=30
Environment=AIPC_LLM_VOICE_MAX_RETRIES=0
