# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-classifier-model.conf
[Service]
# Intent classify on qwythos-9b (small Mythos) — uncensored routing + warm speed.
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER_TIMEOUT=3.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
