# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-uncensored-agent.conf
[Service]
# Hybrid 2026-07-10 (speed + no-filter routing/chat):
#
# qwythos-9b (~5.6GB Mythos, Vulkan) — chat + intent classify
#   - Warm ≈0.2–0.5s; uncensored routing (not NPU gemma-it)
#   - One model serves both so classify does not thrash a second weight
# coder-agentic (Qwen3.6 Uncensored Aggressive) — Hermes / daily tool LLM
# vlm-screen — fast desktop describe
# resident-small (NPU) — optional max-speed only if you override SUPERVISOR
#
# Override max-speed NPU chat (will re-introduce IT refusals on chat):
#   Environment=AIPC_SUPERVISOR_MODEL=resident-small
Environment=AIPC_SUPERVISOR_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_TIMEOUT=3.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
Environment=AIPC_DAILY_MODEL=coder-agentic
Environment=AIPC_DAILY_DIRECT_TOOLS=1
Environment=AIPC_SCREEN_VLM=vlm-screen
Environment=AIPC_LLM_VOICE_TIMEOUT=25
Environment=AIPC_LLM_VOICE_MAX_RETRIES=0
