# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-uncensored-agent.conf
[Service]
# Agent stack → local uncensored aliases (see models.yaml / litellm).
# Chat + tools + Hermes share coder-agentic (Qwen3.6 Uncensored Aggressive)
# so max_loaded_models=2 is not thrashed by 3 different 20GB weights.
# Classify: qwythos-9b (Mythos, small). Screen VLM: Gemma4 Uncensored Vision.
# resident-small (FLM IT) stays for non-agent closed-loop / NPU pin only.
# Optional softer chat: AIPC_SUPERVISOR_MODEL=assistant-gemma
Environment=AIPC_SUPERVISOR_MODEL=coder-agentic
Environment=AIPC_DAILY_MODEL=coder-agentic
Environment=AIPC_CLASSIFIER_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_SCREEN_VLM=vlm-qwen2vl
