# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-workers.conf
[Service]
# Two workers so long Hermes/sync tool calls do not block free-chat /
# classify on the same process. STM/pending/session are disk-merged.
ExecStart=
ExecStart=/var/lib/aipc-agent/venv/bin/uvicorn aipc_agent.server:app --host 127.0.0.1 --port 4100 --workers 2
