# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-classifier-model.conf
[Service]
# Wins alphabetically over older speed/zz drop-ins if any still set rules.
# Classifier on NPU resident-small (warm FLM) for low-latency routing.
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER_TIMEOUT=2.0
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
