# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-uncensored-agent.conf
[Service]
# Hybrid: NPU speed + Vulkan uncensored tools (2026-07-10).
#
# NPU (resident-small / FLM gemma4-e4b) — always warm, ~1–3s voice turns:
#   AIPC_SUPERVISOR_MODEL, AIPC_CLASSIFIER_MODEL
# Vulkan uncensored (coder-agentic Qwen3.6 Aggressive) — tools only:
#   AIPC_DAILY_MODEL, Hermes CLI (~/.hermes model.default)
# Screen: prefer small VLM (vlm-screen ~7B) for “看桌面”; set
#   AIPC_SCREEN_VLM=vlm-qwen2vl only when you need uncensored vision.
#
# Override examples:
#   AIPC_SUPERVISOR_MODEL=coder-agentic   # full uncensored chat (slower)
#   AIPC_SUPERVISOR_MODEL=assistant-gemma  # balanced uncensored chat
Environment=AIPC_SUPERVISOR_MODEL=resident-small
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_TIMEOUT=2.0
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
Environment=AIPC_DAILY_MODEL=coder-agentic
Environment=AIPC_DAILY_DIRECT_TOOLS=1
Environment=AIPC_SCREEN_VLM=vlm-screen
Environment=AIPC_LLM_VOICE_TIMEOUT=20
Environment=AIPC_LLM_VOICE_MAX_RETRIES=0
