# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-skill-learn.conf
[Service]
# ── Learning stack: ON by default ─────────────────────────────────
# Hermes trail → async mentor → on-box skills → match inject.
# No topic/keyword gates. Set any var to 0/off only to debug.

# Background learn queue (never blocks voice TTS)
Environment=AIPC_LEARN_BG=1
Environment=AIPC_LEARN_QUEUE_MAX=64

# Need-triggered skill extract after successful tool turns
Environment=AIPC_SKILL_LEARN=1
Environment=AIPC_SKILL_LEARN_MODEL=ornith-35b
Environment=AIPC_SKILL_LEARN_FALLBACK_MODEL=assistant-gemma
Environment=AIPC_SKILL_LEARN_TIMEOUT=120

# Local skill tree (never under aipc modules/)
Environment=AIPC_SKILL_ROOT=/var/lib/aipc-agent/skills
Environment=AIPC_EPISODE_DIR=/var/lib/aipc-agent/episodes

# Hermes can use tools + isolated browser (model decides when to call)
Environment=AIPC_HERMES_ROUTE=1
Environment=AIPC_HERMES_BROWSER=auto
# Multi-engine search cold-start for Hermes (searxng/ddg/brave/bing).
# Side-path sites OK. auto|1|0 — default auto = inject when browser on.
Environment=AIPC_WEB_HINT=1
Environment=AIPC_WEB_HINT_HERMES=auto
Environment=AIPC_SEARXNG_ENDPOINT=http://127.0.0.1:8888

# Durable memory of successful turns (feeds later self-improve)
Environment=AIPC_MEM0_INTERNALIZE=1
