# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzz-classifier-model.conf
[Service]
# Wins alphabetically over older speed/zz drop-ins if any still set rules.
# Classifier uses small Mythos local model (uncensored agent stack).
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_MODEL=qwythos-9b
Environment=AIPC_CLASSIFIER_TIMEOUT=2.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
