# modules/agent-orchestrator/files/etc/systemd/system/aipc-agent-orchestrator.service.d/zzzz-smo-memory-guard.conf
[Service]
# SMO memory guard (2026-07-16): background services must NOT cold-load
# mid/large Vulkan models. Chat/classify/learn stay on NPU resident-small.
# Tool/coding path may still use coder-agentic (on the GPU allowlist).
# Heavy 122B is user-selected via Hermes/LiteLLM only.
Environment=AIPC_SUPERVISOR_MODEL=resident-small
Environment=AIPC_CLASSIFIER_MODEL=resident-small
Environment=AIPC_CLASSIFIER=always
Environment=AIPC_CLASSIFIER_TIMEOUT=2.5
Environment=AIPC_CLASSIFIER_MAX_TOKENS=24
Environment=AIPC_DAILY_MODEL=coder-agentic
Environment=AIPC_SKILL_LEARN_MODEL=resident-small
Environment=AIPC_SKILL_LEARN_FALLBACK_MODEL=resident-small
Environment=AIPC_TECH_ADVISOR_MODEL=coder-agentic
# Prefer NPU for free-chat; tools still can hit agentic when needed.
Environment=AIPC_LLM_VOICE_TIMEOUT=20
