# modules/voice-pipecat/files/usr/lib/systemd/user/aipc-bluetooth-audio-recover.service
[Unit]
Description=Recover paired Bluetooth audio after a half-connect
After=pipewire.service pipewire-pulse.service wireplumber.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=oneshot
ExecStart=/etc/aipc/aipc-bluetooth-audio-recover
TimeoutStartSec=60

[Install]
WantedBy=default.target
