# modules/voice-pipecat/files/usr/lib/systemd/user/aipc-bluetooth-audio-recover.service
[Unit]
Description=Watch Bluetooth and recover a half-connected audio speaker
After=pipewire.service pipewire-pulse.service wireplumber.service bluetooth.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=simple
ExecStart=/etc/aipc/aipc-bluetooth-audio-recover --monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
