# modules/agent-orchestrator/files/etc/systemd/system/aipc-self-improve.service
[Unit]
Description=AIPC background self-improve (episode → local skills)
Documentation=file:///var/home/birdyo/aipc-strix-halo/openspec/changes/assistant-self-improvement/design.md
After=network-online.target aipc-agent-orchestrator.service
Wants=network-online.target

[Service]
Type=oneshot
# Live hotfix path (ostree: process lives under /var/lib/aipc-agent)
Environment=AIPC_EPISODE_DIR=/var/lib/aipc-agent/episodes
Environment=AIPC_SKILL_ROOT=/var/lib/aipc-agent/skills
Environment=AIPC_LEARN_STATE=/var/lib/aipc-agent/learning/self_improve_state.json
# Learning mentor ON by default (same as orchestrator)
Environment=AIPC_SKILL_LEARN=1
Environment=AIPC_SKILL_LEARN_MODEL=ornith-35b
Environment=AIPC_SKILL_LEARN_FALLBACK_MODEL=assistant-gemma
Environment=AIPC_SKILL_LEARN_TIMEOUT=120
Environment=AIPC_SELF_IMPROVE_MAX=15
Environment=AIPC_LITELLM_URL=http://127.0.0.1:4000
# Nice/idle-ish: do not contend with voice critical path
Nice=15
IOSchedulingClass=idle
Environment=PYTHONPATH=/var/lib/aipc-agent
WorkingDirectory=/var/lib/aipc-agent
ExecStart=/var/lib/aipc-agent/venv/bin/python -m aipc_agent self-improve --hours 72 --max 15
# Allow long LLM extract batch
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
