# modules/agent-orchestrator/files/etc/systemd/system/aipc-self-improve.timer
[Unit]
Description=AIPC background self-improve timer (idle skill growth)
Documentation=file:///var/home/birdyo/aipc-strix-halo/openspec/changes/assistant-self-improvement/design.md

[Timer]
# Nightly + a few daytime idle passes
OnCalendar=*-*-* 03:30:00
OnCalendar=*-*-* 14:00:00
# No OnBootSec: a post-boot pass thrashed Lemonade on an idle machine (0011).
# The preflight in self_improve.py additionally skips busy/battery windows.
Persistent=true
RandomizedDelaySec=10min
Unit=aipc-self-improve.service

[Install]
WantedBy=timers.target
