# modules/voice-pipecat/quadlet/pipecat.container
[Unit]
Description=AIPC Pipecat voice orchestrator
After=litellm.service voice-wake.service sensevoice.service paraformer.service cosyvoice.service kokoro.service
Requires=litellm.service

[Container]
Image=docker.io/pipecat/pipecat:latest
ContainerName=aipc-voice-pipecat
EnvironmentFile=/etc/aipc/env.d/voice-pipecat/config
Volume=/etc/aipc/pipecat:/etc/aipc/pipecat:ro,Z

[Service]
Restart=always

[Install]
WantedBy=default.target
