# modules/voice-wake/files/etc/systemd/system/aipc-voice-wake.service
[Unit]
Description=aipc voice wake (custom phrase / openWakeWord / energy)
After=sound.target sleep.target user@1000.service
Requires=user@1000.service
Conflicts=aipc-voice-mute.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=birdyo
Group=birdyo
Environment=HOME=/var/home/birdyo
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=PYTHONPATH=/var/lib/aipc-voice/lib
Environment=AIPC_VOICE_ONCE=/var/home/birdyo/.local/bin/aipc-voice-once
Environment=AIPC_VOICE_MUTE_FLAG=/run/aipc/voice-mute
Environment=AIPC_WAKE_PHRASES_FILE=/etc/aipc/wake/phrases
Environment=AIPC_VOICE_STT_URL=http://127.0.0.1:9001/transcribe
Environment=AIPC_WAKE_POLICY_FILE=/etc/aipc/voice/wake-policy.env
EnvironmentFile=-/etc/aipc/voice/wake-policy.env
ExecStart=/usr/bin/python3 /var/lib/aipc-voice/lib/aipc_voice_wake.py --mode phrase
Restart=always
RestartSec=10
ExecStartPre=+/bin/mkdir -p /var/lib/aipc-voice/wake /run/aipc /etc/aipc/wake /etc/aipc/voice
ExecStartPre=+/bin/chmod 755 /run/aipc

[Install]
WantedBy=multi-user.target
