# modules/voice-wake/files/etc/systemd/system/aipc-ai-sleep-pre.service
[Unit]
Description=AIPC AI voice stack — prepare for system sleep
DefaultDependencies=no
Before=sleep.target
Conflicts=umount.target

[Service]
Type=oneshot
ExecStart=/var/lib/aipc-voice/bin/aipc-ai-stack-sleep pre suspend
TimeoutStartSec=15

[Install]
WantedBy=sleep.target
