# modules/llm-lemonade/files/etc/systemd/system/aipc-lemonade-idle-release.timer
[Unit]
Description=Periodic Lemonade idle release for compact aliases
Documentation=file:///var/home/birdyo/aipc-strix-halo/openspec/changes/0006-lemonade-compact-idle-unload/specs/ai-runtime/spec.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=aipc-lemonade-idle-release.service

[Install]
WantedBy=timers.target
