# modules/agent-screen-control/files/etc/systemd/system/ydotool.service.d/override.conf
# ydotool.service override — agent-screen-control (phase-4-agent#4.7)
#
# The base unit (from the ydotool package) is `ExecStart=/usr/bin/ydotoold`
# with no args, which uses ydotoold's default (root-only) socket. This
# override points it at a uid-1000-owned socket path instead, so input.py's
# unprivileged `ydotool` calls can open it.
#
# Root-cause fix for a boot-time race found live on 2026-07-11: this socket
# path lives under /run/user/1000, which is created by
# user-runtime-dir@1000.service — not guaranteed to exist yet when
# ydotool.service (WantedBy=default.target) starts early in boot. Without an
# explicit dependency, ydotoold sometimes starts before that directory
# exists, fails to bind the socket, hits Restart=always's start-limit burst,
# and gives up (status=2, ydotool.service ends up failed) until something
# runs `systemctl restart ydotool` by hand — which always succeeds because
# by then the directory exists. Restart=always alone doesn't fix this: the
# burst limit is hit faster than the directory shows up.
#
# The After=/Requires= pair below orders ydotool.service behind the runtime
# dir unit and pulls it in as a dependency, so /run/user/1000 always exists
# before ydotoold tries to bind — removes the race instead of retrying
# around it.
[Unit]
After=user-runtime-dir@1000.service
Requires=user-runtime-dir@1000.service

[Service]
ExecStart=
ExecStart=/usr/bin/ydotoold --socket-path=/run/user/1000/.ydotool_socket --socket-own=1000:1000
