# modules/system-unified-memory/files/etc/systemd/system/gpp-wake-fix.service
[Unit]
Description=Disable spurious GPP0/GPP1 ACPI wake sources (GZ302EA sleep/wake fix)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/aipc/gpp-wake-fix

[Install]
WantedBy=default.target
