# modules/system-unified-memory/files/etc/systemd/system/platform-profile-idle-check.service
[Unit]
Description=Check idle CPU/GPU load and adjust ASUS EC platform power profile accordingly (GZ302EA)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/aipc/platform-profile-idle-check
