# modules/system-unified-memory/files/etc/systemd/system/platform-profile-auto.service
[Unit]
Description=Apply ASUS EC platform power profile based on AC adapter state (GZ302EA)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/aipc/platform-profile-auto

[Install]
WantedBy=default.target
