# modules/system-unified-memory/files/etc/systemd/system/platform-profile-idle-check.timer
[Unit]
Description=Periodic idle check for ASUS EC platform power profile (GZ302EA)

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
