# modules/voice-wake/files/etc/systemd/system/aipc-voice-wake.service.d/zzzzzzz-anti-ghost-strict.conf
[Service]
# 2026-07-16: ambient STT 我。 + score~78 was fuzzy_promoted — prefer miss.
Environment=AIPC_WAKE_ALLOW_FUZZY_PROMOTE=0
Environment=AIPC_WAKE_PROMOTE_SCORE=90
Environment=AIPC_WAKE_CANDIDATE_SCORE=70
Environment=AIPC_WAKE_ENERGY=2800
Environment=AIPC_WAKE_ENERGY_FRAMES=8
Environment=AIPC_WAKE_COOLDOWN=10
