# modules/voice-pipecat/files/etc/systemd/system/aipc-voice-mute.target
[Unit]
Description=AIPC voice mute target
