# modules/voice-tts-kokoro/quadlet/kokoro.container
[Unit]
Description=AIPC Kokoro TTS (en fallback)

[Container]
Image=docker.io/aipc/kokoro:latest
ContainerName=aipc-kokoro
PublishPort=127.0.0.1:9012:9012

[Service]
Restart=always

[Install]
WantedBy=default.target
