# modules/memory-mem0/files/etc/systemd/system/aipc-mem0-state-dir.service
[Unit]
Description=Prepare aipc-mem0 state directory ownership
Before=aipc-mem0.service

[Service]
Type=oneshot
ExecStart=/usr/lib/aipc-mem0/aipc-mem0-state-dir-setup
RemainAfterExit=yes

[Install]
WantedBy=default.target
