# modules/dev-ai-hermes-webui/files/etc/systemd/system/aipc-hermes-webui-setup.service
[Unit]
Description=Provision Hermes WebUI home checkout + linger (first boot / pin change)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/lib/aipc/hermes-webui/setup-hermes-webui.sh

[Service]
Type=oneshot
ExecStart=/usr/lib/aipc/hermes-webui/setup-hermes-webui.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
