# modules/rag-ingest/files/etc/systemd/system/aipc-rag-screen-audio.service
[Unit]
Description=AIPC RAG watcher — screen OCR + audio transcript (opt-in, TTL-gated)
After=network-online.target rag-embedder.service postgres.service
Requires=rag-embedder.service postgres.service

[Service]
Type=simple
Environment=PYTHONPATH=/usr/lib/aipc-rag
ExecStart=/usr/lib/aipc-rag/venv/bin/python3 -m aipc_rag.screen_audio
Restart=on-failure

[Install]
WantedBy=default.target
